# xmpd.service
[Unit]
Description=Multi-source music daemon (xmpd)
Documentation=https://github.com/tuncenator/xmpd
After=network.target

[Service]
Type=simple
# Update the path below to match your installation directory
ExecStart=/path/to/xmpd/.venv/bin/python -m xmpd
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.config/xmpd %h/Music

[Install]
WantedBy=default.target
